# ansible/roles/redis/templates/ares@.service.j2
[Unit]
Description=Ares Worker (%i)
After=redis.service nats-server.service
Wants=redis.service nats-server.service

[Service]
Type=simple
ExecStart={{ redis_ares_worker_binary }} worker
Environment=ARES_REDIS_URL=redis://{{ redis_bind_address }}:{{ redis_port }}
Environment=NATS_URL=nats://{{ redis_bind_address }}:4222
Environment=ARES_WORKER_ROLE=%i
Environment=ARES_WORKER_MODE=tool_exec
Environment=RUST_LOG=info
Restart=on-failure
RestartSec=5
StandardOutput=append:{{ redis_ares_log_dir }}/%i.log
StandardError=append:{{ redis_ares_log_dir }}/%i.log

# Contain child processes (netexec, hashcat, nmap, etc.) within this cgroup.
# Without these limits, runaway tool processes can OOM the entire system and
# take down the SSM agent.
Delegate=yes
Slice=system-ares.slice
MemoryHigh={{ redis_ares_worker_memory_high }}
MemoryMax={{ redis_ares_worker_memory_max }}
TasksMax={{ redis_ares_worker_tasks_max }}

[Install]
WantedBy=multi-user.target
